# keys/mfsbsdonly/etc/rc.d/jautogui
#!/bin/sh
# Copyright (c) 2014 John Ko
# $Id$

# PROVIDE: jautogui
# REQUIRE: mountcritlocal disks jautomount mdinit tmp var jroute jautopkg
# KEYWORD: FreeBSD

. /etc/rc.subr

name="jautogui"
start_cmd="jautogui_start"
stop_cmd=":"

jautogui_start()
{
	autogui=`/bin/kenv -q autogui`
	case ${autogui} in
		[yY][eE][sS])
			test -x /server/bin/startslim && /server/bin/startslim
			;;
	esac
}

load_rc_config $name
run_rc_command "$1"
